# categories/post-install tweaks/13. hyprland/.config/systemd/user/hyprland-kded6-sni-disable.service
[Unit]
Description=Prevent kded6 from stealing statusnotifierwatcher
PartOf=hyprland-session.target
After=dbus.socket

[Service]
Type=simple
ExecStart=%h/.local/bin/sni-guard

[Install]
WantedBy=hyprland-session.target
