# categories/post-install tweaks/13. hyprland/.config/systemd/user/hyprland-kded6-sni-disable.service
[Unit]
Description=Prevent kded6 from stealing statusnotifierwatcher
After=dbus.socket hyprland-session-ready.service
PartOf=graphical-session.target
Before=dms.service

[Service]
ExecStart=%h/.local/bin/sni-guard
Restart=on-failure

[Install]
WantedBy=hyprland-session-ready.service
