# categories/post-install tweaks/13. hyprland/etc/systemd/user/hyprland-session-ready.service
[Unit]
Description=Signal Hyprland session readiness and propagate environment
After=hyprland.service graphical-session.target
PartOf=hyprland-session.target
BindsTo=graphical-session.target
Wants=graphical-session.target xdg-desktop-autostart.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/hyprland-session-ready
