# categories/post-install tweaks/13. hyprland/etc/systemd/user/hyprland-session.target
[Unit]
Description=Hyprland compositor session
Wants=graphical-session-pre.target
After=graphical-session-pre.target
